# conf/supervisord.conf
[supervisord]
nodaemon=true

[program:app]
directory=/app/python/
command=python3 promptoon.py
autostart=true
autorestart=true

